# deploy/crontab
# Crontab de supercronic — reemplaza la Tarea Programada de Windows (instalar-tarea.ps1)
# Equivalente: Lunes a Viernes, cada 5 minutos, de 6:00am a 6:55pm.
# La hora se evalúa en la zona horaria del contenedor (TZ=America/Bogota, ver docker-compose.yml).
# Por defecto supercronic no solapa ejecuciones: si "node index.js" sigue corriendo cuando
# toca el siguiente disparo, ese disparo se omite (igual que "MultipleInstances IgnoreNew" en Windows).
*/5 6-18 * * 1-5 cd /app && node index.js
